# docker/supervisor/conf.d/rabbitmq.conf
[program:rabbitmq]
command=/usr/sbin/rabbitmq-server
numprocs=1
autostart=true
autorestart=true

